# rate-mirrors.timer
[Unit]
Description=Run rate-mirrors

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=3h

[Install]
WantedBy=timers.target
